# remote_bridge/omniremote-bridge.service
[Unit]
Description=OmniRemote USB/Bluetooth Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /opt/omniremote-bridge/omniremote_bridge.py -c /etc/omniremote-bridge/config.json
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
